# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command= /etc/init.d/ssh start -D


[program:nginx]
command= /usr/local/nginx/sbin/nginx 

